# etc/crontab
*/2 *   * * *   root    /opt/gsm/generate_points.sh 2>/dev/null
* *     * * *   root    /opt/gsm/monitoring_services.sh 2>/dev/null
